# launchd/launched.testNew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>launched.test_timer</string>

    <!-- Use the absolute path to the script; no need for sh -c -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ff_tc_25006/github/project-backup/testNew.sh</string>
    </array>

    <!-- Run every 20 seconds (or whatever you need) -->
    <key>StartInterval</key>
    <integer>20</integer>

    <!-- Optional: run once at load time -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Capture launchd’s own stdout/stderr (useful for debugging) -->
    <key>StandardOutPath</key>
    <string>/Users/ff_tc_25006/Library/Logs/launched.test_timer.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ff_tc_25006/Library/Logs/launched.test_timer.err.log</string>
</dict>
</plist>
